# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock              ; (the path to the socket file)

[supervisord]
logfile=%(here)s/logs/supervisord.log  ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB                  ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10                     ; (num of main logfile rotation backups;default 10)
loglevel=info                          ; (log level;default info; others: debug,warn,trace)
pidfile=%(here)s/supervisord.pid       ; (supervisord pidfile;default supervisord.pid)
nodaemon=false                         ; (start in foreground if true;default false)
minfds=1024                            ; (min. avail startup file descriptors;default 1024)
minprocs=200                           ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock                          ; use a unix:// URL  for a unix socket

[program:server]
command=gunicorn -b localhost:8000 -c gunicorn.conf manage:app ; the program (relative uses PATH, can take args)
;process_name=%(program_name)s                                 ; process_name expr (default %(program_name)s)
;directory=/tmp                                                ; directory to cwd to before exec (def no cwd)
numprocs=1                                                     ; number of processes copies to start (def 1)
autostart=true                                                 ; start at supervisord start (default: true)
autorestart=unexpected                                         ; whether/when to restart (default: unexpected)
;startretries=3                                                ; max # of serial start failures (default 3)
startsecs=0                                                    ; number of secs prog must stay running (def. 1)
stopwaitsecs=0                                                 ; max num secs to wait b4 SIGKILL (default 10)
stdout_logfile=%(here)s/logs/%(program_name)s.log.out          ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB                                   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile=%(here)s/logs/%(program_name)s.log.err          ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=50MB                                   ; max # logfile bytes b4 rotation (default 50MB)
